# .config/feh/time-backgrounds/background-changer.service
[Unit]
Description=Background Changer Service

[Service]
ExecStart=sh %h/.config/feh/time-backgrounds/background-changer.sh

