# .config/feh/time-backgrounds/background-changer.timer
[Unit]
Description=Change Background Every 2 hours

[Timer]
OnBootSec=0
OnCalendar=00/2:00
Persistent=true

[Install]
WantedBy=timers.target
