# .config/feh/time-backgrounds/background-changer.service
[Unit]
Description=Background Changer Service

[Service]
ExecStart=/bin/bash %h/.config/feh/time-backgrounds/background-changer.sh

